# assets/etc/supervisor/conf.d/php-fmp.ini
[program:php-fpm]
command=/bin/bash -c "mkdir -p /var/run/php ; \
                      /usr/sbin/php-fpm7.2 -F -c /etc/php/7.2/fpm/php.ini --fpm-config /etc/php/7.2/fpm/php-fpm.conf"

autostart=true
autorestart=true
stopasgroup=true

